# lcdctl.service
[Unit]
Description=QNAP TS-853A LCD Controller (stats mode)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/lcdctl stats --interval 5
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
